# root/etc/periodic/daily/tabby-cron
#!/bin/bash
# Run the cron script for Tabby

php /config/www/tabby/cron.php